# deploy/systemd/actap-frontend.service
# ACTAP frontend (nginx) systemd unit
# Copy to /etc/systemd/system/actap-frontend.service

[Unit]
Description=ACTAP Frontend (Nginx serving built static assets)
After=network-online.target actap-backend.service
Wants=network-online.target

[Service]
Type=forking
PIDFile=/run/nginx.pid
ExecStartPre=/usr/sbin/nginx -t
ExecStart=/usr/sbin/nginx
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s QUIT $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10
Restart=always
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
